# debian/coda-server.conffiles
/etc/init.d/auth2.init
/etc/init.d/codasrv.init
