# coda-src/scripts/venus.init.in
#!/bin/sh
#
# skeleton      Example file to build /etc/init.d scripts.
#
# Version:      @(#) /etc/init.d/skeleton 1.01 26-Oct-1993
#
# Author:       Miquel van Smoorenburg, <miquels@drinkel.nl.mugnet.org>
#

# chkconfig: 345 95 01
# description: The Coda cache manager.


PATH=/sbin:/usr/sbin:/usr/local/sbin:$PATH
export PATH

# defaults
pid_file=/usr/coda/venus.cache/pid

# overridden from the venus.conf
[ -f @sysconfdir@/venus.conf ] && . @sysconfdir@/coda/venus.conf

[ -f @sbindir@/venus ] || exit 0

echon() {
    if [ "`echo -n`" ] ; then
        echo "$@"\c
    else
        echo -n "$@"
    fi
}


# See how we were called.
case "$1" in
  start)
	echon "Starting venus: "
	if [ -f $pid_file ]; then
		echo "$pid_file exists, venus already running?."
		exit 1
	fi
	@sbindir@/venus &
	echo "done."
    ;;
  stop)
	echon "Shutting down venus: "
	if [ -f $pid_file ]; then
		kill -TERM `cat $pid_file`
		umount /coda
	fi
	rm -f $pid_file
	echo "done."
    ;;
  hardstop)
	echon "Killing venus: "
	if [ -f $pid_file ]; then
		kill -9 `cat $pid_file`
		umount /coda
	fi
	rm -f $pid_file
	echo "done."
    ;;
    restart|force-reload)
	echo "$1 not available"
    ;;
    *)
	echo "Usage: $0 {start|stop|hardstop}"
	exit 1
    ;;
esac

exit 0
